# github_backup/system/launchagents/ai.openclaw.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.gateway</string>
    
    <key>Comment</key>
    <string>OpenClaw Gateway (v2026.4.15)</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>1</integer>
    <key>Umask</key>
    <integer>63</integer>
    <key>ProgramArguments</key><array><string>/bin/bash</string><string>/Users/KS/.agents/skills/openclaw-fortress-1.0.0/logic/start_gateway.sh</string></array>
    
    <key>StandardOutPath</key>
    <string>/Users/KS/.openclaw/logs/gateway.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/KS/.openclaw/logs/gateway.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
    <key>NODE_OPTIONS</key>
    <string>--import /Users/KS/.openclaw/guardian/epperm-patch.mjs</string>
    <key>NODE_USE_SYSTEM_CA</key>
    <string>1</string>
    <key>GEMINI_API_KEY</key><string>[GCP_API_KEY_REDACTED]</string>
    <key>GOOGLE_API_KEY</key><string>[GCP_API_KEY_REDACTED]</string>
    <key>GOOGLE_AI_API_KEY</key><string>[GCP_API_KEY_REDACTED]</string>
    <key>HOME</key>
    <string>/Users/KS</string>
    <key>TMPDIR</key>
    <string>/var/folders/nj/2ptqq9p16ng0mtv_5qt7y1300000gq/T/</string>
    <key>NODE_EXTRA_CA_CERTS</key>
    <string>/etc/ssl/cert.pem</string>
    <key>PATH</key>
    <string>/usr/sbin:/sbin:/usr/bin:/bin:/opt/homebrew/opt/node/bin:/Users/KS/.local/bin:/Users/KS/.npm-global/bin:/Users/KS/bin:/Users/KS/.volta/bin:/Users/KS/.asdf/shims:/Users/KS/.bun/bin:/Users/KS/Library/Application Support/fnm/aliases/default/bin:/Users/KS/.fnm/aliases/default/bin:/Users/KS/Library/pnpm:/Users/KS/.local/share/pnpm:/opt/homebrew/bin:/usr/local/bin</string>
    <key>OPENCLAW_GATEWAY_PORT</key>
    <string>18792</string>
    <key>OPENCLAW_LAUNCHD_LABEL</key>
    <string>ai.openclaw.gateway</string>
    <key>OPENCLAW_SYSTEMD_UNIT</key>
    <string>openclaw-gateway.service</string>
    <key>OPENCLAW_WINDOWS_TASK_NAME</key>
    <string>OpenClaw Gateway</string>
    <key>OPENCLAW_SERVICE_MARKER</key>
    <string>openclaw</string>
    <key>OPENCLAW_SERVICE_KIND</key>
    <string>gateway</string>
    <key>OPENCLAW_SERVICE_VERSION</key>
    <string>2026.4.15</string>
    <key>[SENSITIVE_TOKEN_HARD_REDACTED]</key>
    <string>GEMINI_API_KEY,GOOGLE_AI_API_KEY,GOOGLE_API_KEY</string>
    </dict>
  </dict>
</plist>
