# github_backup/system/launchagents/ai.openclaw.daily.optimizer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.daily.optimizer</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/KS/.openclaw/workspace/config/governance/openclaw_daily_optimizer.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/KS/.openclaw/workspace/logs/daily_optimizer.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/KS/.openclaw/workspace/logs/daily_optimizer.err.log</string>
</dict>
</plist>
