# github_backup/system/launchagents/ai.openclaw.guardian.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.guardian</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/KS/.openclaw/guardian/guardian.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/KS/.openclaw/logs/system-guardian.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/KS/.openclaw/logs/system-guardian.launchd.err.log</string>
</dict>
</plist>
