# github_backup/system/launchagents/ai.openclaw.soul-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.soul-backup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/KS/.agents/skills/openclaw-soul-sync-2.3.1/soul_sync.py</string>
        <string>--mode</string>
        <string>full</string>
        <string>--github</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>GITHUB_TOKEN</key>
        <string>gho_[GENERIC_SECRET_REDACTED]</string>
        <key>SOUL_SYNC_REPO</key>
        <string>scchin/openclaw-soul-backups</string>
        <key>HOME</key>
        <string>/Users/KS</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/KS/.openclaw/logs/soul_backup.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/KS/.openclaw/logs/soul_backup.err</string>
</dict>
</plist>
